# server-py/src/flask_app/supervisor/conf.d/gunicorn.conf
[program:gunicorn]
command=/code/lib/batch/app.sh
directory=/code
user=nobody
autostart=true
autorestart=true
redirect_stderr=true